# deploy/systemd/agent-suite-backup.timer
# agent-suite-backup.timer — nightly trigger for the backup+verify-restore service
# Generated by `agent-suite schedule install`. Do not edit by hand.

[Unit]
Description=Nightly pg_dump of the suite Postgres store + weekly verify-restore (timer)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
