# deploy/systemd/agent-suite-doctor-alert.timer
# agent-suite-doctor-alert.timer — hourly trigger for the doctor+alert service
# Generated by `agent-suite schedule install`. Do not edit by hand.

[Unit]
Description=Periodic doctor health check + alert routing on state change (timer)

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
