# deploy/systemd/agent-suite-doctor-alert.service
# agent-suite-doctor-alert.service — periodic doctor health check + alert routing
# Generated by `agent-suite schedule install`. Do not edit by hand.

[Unit]
Description=Periodic doctor health check + alert routing on state change
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
ExecStart=agent-suite alert-check
EnvironmentFile=-/etc/agent-suite/suite.env
User=root

[Install]
WantedBy=multi-user.target
