# deploy/systemd/agent-suite-backup.service
# agent-suite-backup.service — nightly pg_dump + verify-restore
# Generated by `agent-suite schedule install`. Do not edit by hand.
# Re-run `agent-suite schedule install` to regenerate.

[Unit]
Description=Nightly pg_dump of the suite Postgres store + weekly verify-restore
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
ExecStart=agent-suite backup --verify-restore
EnvironmentFile=-/etc/agent-suite/suite.env
User=root

[Install]
WantedBy=multi-user.target
